# dashboard.service
[Unit]
Description=System Health Dashboard

[Service]
Type=oneshot
User=yourusername
ExecStart=/bin/bash /home/yourusername/dashboard.sh
StandardOutput=append:/home/yourusername/dashboard.log
StandardError=append:/home/yourusername/dashboard.log
